# debian/plegger.postinst
#!/bin/sh
set -e if [ "$1" = 'configure' -a -n "$2" ]; then
    # install
    echo "DAEMON_CONF=\"/etc/hostapd/hostapd.conf\"" > /etc/default/hostapd
    systemctl enable dnsmasq
    systemctl enable hostapd
    exit 0

#DEBHELPER#
